# service/tor/run
#!/bin/sh

/sbin/su-exec tor /usr/bin/tor -f /etc/service/tor/config
