# service/privoxy/run
#!/bin/sh

/sbin/su-exec privoxy /usr/sbin/privoxy --no-daemon
